# axiom-validator.service
[Unit]
Description=AXIOM Protocol Validator Node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=codespace
WorkingDirectory=/workspaces/Qubit-Protocol-84m
ExecStart=/workspaces/Qubit-Protocol-84m/target/release/qubit --config /workspaces/Qubit-Protocol-84m/axiom-validator.toml
Restart=always
RestartSec=10
StandardOutput=append:/workspaces/Qubit-Protocol-84m/./axiom-validator-data/logs/validator.log
StandardError=append:/workspaces/Qubit-Protocol-84m/./axiom-validator-data/logs/validator-error.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only

# Resource limits
LimitNOFILE=65536
MemoryMax=16G

[Install]
WantedBy=multi-user.target
